# packages/versionchecker/fastboot/systemd/system/versionchecker.timer
[Unit] 
Description=Runs versionchecker.sh every 10 minuntes

[Timer] 
OnBootSec=10min 
OnCalendar=*:0/10
Unit=versionchecker.service 

[Install] 
WantedBy=multi-user.target