# packages/versionchecker/etc/systemd/system/versionchecker.service
[Unit]
Description=versionchecker custom service

[Service]
Type=simple
ExecStart=/bin/versionchecker.sh
ExecReload=/usr/bin/versionchecker.sh
[Install]
WantedBy=multi-user.target