# packages/versionchecker/etc/systemd/system/versionchecker.timer
[Unit] 
Description=Runs versionchecker.sh on boot after 10 minutes then every 30 minutes.

[Timer] 
OnBootSec=10min
OnUnitActiveSec=30min 
 