# packages/versionchecker/etc/systemd/system/versionchecker.service
[Unit]
Description=versionchecker custom service

[Service]
Type=simple
EnvironmentFile=/etc/thinstation.env
ExecStart=/bin/versionchecker.sh
ExecReload=/bin/versionchecker.sh
SyslogIdentifier=thinstation

[Install]
WantedBy=versionchecker.time