# packages/versionchecker/etc/systemd/system/versionchecker.service
[Unit]
Description=versionchecker custom service

[Service]
Type=idle
EnvironmentFile=/etc/thinstation.env
ExecStart=/bin/versionchecker.sh
SyslogIdentifier=thinstation
