# packages/versionchecker/etc/systemd/system/versionchecker.timer
[Unit] 
Description=Runs versionchecker.sh every 10 minuntes

[Timer] 
OnBootSec=5min
OnUnitActiveSec=10min 
 