# packages/assetreporter/etc/systemd/system/assetreporter.service
[Unit]
Description=assetreporter custom service

[Service]
Type=idle
ExecStart=/bin/assetreporter.sh
SyslogIdentifier=thinstation
