# packages/assetreporter/etc/systemd/system/assetreporter.timer
[Unit] 
Description=Runs assetreporter.sh on boot after 10 minutes then daily.

[Timer] 
OnBootSec=10min
OnUnitActiveSec=daily
RandomizedDelaySec=60