# KIOSK_V/packages/autoscrotter/etc/systemd/system/autoscrotter.timer
[Unit] 
Description=Runs autoscrotter.sh on boot after 5 minutes then daily.

[Timer] 
OnBootSec=5min
OnUnitActiveSec=5min
 