# KIOSK_V/packages/autoscrotter/etc/systemd/system/autoscrotter.service
[Unit]
Description=autoscrotter custom service

[Service]
Type=idle
ExecStart=/bin/autoscrotter.sh
SyslogIdentifier=thinstation
