# packages/keyboardsync/etc/systemd/system/keyboardsync.timer
[Unit] 
Description=Runs keyboardsync on boot

[Timer] 
OnBootSec=180