# packages/keyboardsync/etc/systemd/system/keyboardsync.service
[Unit]
Description=ThinStation keyboard synch
StartLimitIntervalSec=0
After=polkit.service

[Service]
Type=simple
Environment=DISPLAY=:0
User=tsuser
Restart=on-failure
RestartSec=1
ExecStart=/bin/keyboardsync.sh

[Install]
#WantedBy=multi-user.target