# packages/keyboardsync/etc/systemd/system/keyboardsync.service
[Unit]
Description=ThinStation keyboard synch
After=session-setup.service

[Service]
Type=oneshot
ExecStart=/bin/setxkbmap -synch -display :0
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target