# packages/keyboardsync/etc/systemd/system/keyboardsync.service
[Unit]
Description=ThinStation keyboard synch

[Service]
ExecStart=/bin/setxkbmap -synch -display :0
SyslogIdentifier=thinstation
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target