# packages/keyboardsync/etc/systemd/system/keyboardsync.service
[Unit]
Description=ThinStation keyboard synch

[Service]
Type=simple
User=tsuser
Restart=on-failure
ExecStart=/bin/keyboardsync.sh

[Install]
WantedBy=multi-user.target