# packages/keyboardsync/etc/systemd/system/keyboardsync.service
[Unit]
Description=ThinStation keyboard synch
StartLimitIntervalSec=0

[Service]
Environment=DISPLAY=:0
User=tsuser
Restart=on-failure
RestartSec=10
ExecStart=/bin/keyboardsync.sh

[Install]
WantedBy=graphical.target