# packages/keyboardsync/etc/systemd/system/keyboardsync.service
[Unit]
Description=ThinStation keyboard synch
After=session-setup.service
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
ExecStart=/bin/setxkbmap -synch -display :0
SyslogIdentifier=thinstation
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target