# packages/keyboardsync/etc/systemd/system/keyboardsync.service
[Unit]
Description=ThinStation keyboard synch
After=profile-setup.service
StartLimitIntervalSec=600
StartLimitBurst=30

[Service]
ExecStart=/bin/setxkbmap -synch -display :0
SyslogIdentifier=thinstation
Restart=on-failure
RestartSec=4

[Install]
WantedBy=multi-user.target