# packages/keyboardsync/etc/systemd/system/keyboardsync.service
[Unit]
Description=ThinStation keyboard synch
StartLimitIntervalSec=0

[Service]
Type=idle
Restart=on-failure
RestartSec=10
ExecStart=/bin/keyboardsync.sh

[Install]
WantedBy=multi-user.target